# tools/miniforge3/etc/profile.d/conda.csh
setenv CONDA_EXE "/Users/severinobione/OMNISEEN_V3/tools/miniforge3/bin/conda";
setenv _CONDA_EXE "/Users/severinobione/OMNISEEN_V3/tools/miniforge3/bin/conda";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "/Users/severinobione/OMNISEEN_V3/tools/miniforge3/bin/python";
setenv _CONDA_ROOT "/Users/severinobione/OMNISEEN_V3/tools/miniforge3";

source "/Users/severinobione/OMNISEEN_V3/tools/miniforge3/lib/python3.13/site-packages/conda/shell/etc/profile.d/conda.csh"